# deb/printer/chipcie-print-server_1.0-1/etc/systemd/system/cups-lpd@.service
[Unit]
Description=CUPS LPD Server
Documentation=man:cups-lpd(8)

[Service]
ExecStart=/usr/lib/cups/daemon/cups-lpd
StandardInput=socket

[Install]
WantedBy=multi-user.target
