# deb/reboot-loop/chipcie-reboot-loop_1.0-1/usr/lib/systemd/system/chipcie-reboot-loop.service
[Unit]
Description=Reboots unit instantly

[Service]
ExecStart=/bin/systemctl reboot 

[Install]
WantedBy=multi-user.target

