# deb/hostname/chipcie-hostname_1.0-1/lib/systemd/system/chipcie-hostname.service
[Unit]
Description=Set hostname based on the id from DOMjudge
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sethostname

[Install]
WantedBy=multi-user.target
